# deploy/examples/custodia-signer.service
# Copyright (c) 2026 Marco Fortina
# SPDX-License-Identifier: AGPL-3.0-only
#
# This file is part of Custodia.
# Custodia is distributed under the GNU Affero General Public License v3.0.
# See the accompanying LICENSE file for details.

[Unit]
Description=Custodia Certificate Signer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=custodia
Group=custodia
Environment=CUSTODIA_SIGNER_ADDR=:9444
Environment=CUSTODIA_SIGNER_TLS_CERT_FILE=/etc/custodia/server.crt
Environment=CUSTODIA_SIGNER_TLS_KEY_FILE=/etc/custodia/server.key
Environment=CUSTODIA_SIGNER_CLIENT_CA_FILE=/etc/custodia/client-ca.crt
Environment=CUSTODIA_SIGNER_ADMIN_SUBJECTS=admin
Environment=CUSTODIA_SIGNER_KEY_PROVIDER=file
Environment=CUSTODIA_SIGNER_CA_CERT_FILE=/etc/custodia/ca.crt
Environment=CUSTODIA_SIGNER_CA_KEY_FILE=/etc/custodia/ca.key
Environment=CUSTODIA_SIGNER_CA_KEY_PASSPHRASE_FILE=/etc/custodia/ca.pass
Environment=CUSTODIA_SIGNER_CRL_FILE=/etc/custodia/client.crl.pem
Environment=CUSTODIA_SIGNER_AUDIT_LOG_FILE=/var/log/custodia/signer-audit.jsonl
ExecStart=/usr/local/bin/custodia-signer
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/etc/custodia
ReadWritePaths=/var/log/custodia
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
